# roles/firezone/templates/firezone.service.j2
[Unit]
Description=Firezone VPN Server
After=docker.service
Requires=docker.service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory={{ fz_firezone_dir }}
ExecStartPre=-/bin/bash -c "test -f {{ fz_firezone_dir }}/docker-compose.yml || (echo 'docker-compose.yml not found in {{ fz_firezone_dir }}'; ls -la {{ fz_firezone_dir }}; exit 1)"
ExecStartPre=-/bin/bash -c "which docker-compose || which docker || (echo 'Neither docker-compose nor docker found'; exit 1)"
ExecStartPre=-/bin/bash -c "docker compose version >/dev/null 2>&1 && echo 'Using docker compose' && docker compose down || (echo 'Using docker-compose' && /usr/local/bin/docker-compose down) || echo 'Compose down failed'"
ExecStartPre=-/bin/bash -c "docker compose version >/dev/null 2>&1 && echo 'Using docker compose' && docker compose pull || (echo 'Using docker-compose' && /usr/local/bin/docker-compose pull) || echo 'Compose pull failed'"
ExecStart=-/bin/bash -c "docker compose version >/dev/null 2>&1 && echo 'Using docker compose' && docker compose up || (echo 'Using docker-compose' && /usr/local/bin/docker-compose up) || (echo 'No compose command worked'; exit 1)"
ExecStop=-/bin/bash -c "docker compose version >/dev/null 2>&1 && docker compose down || /usr/local/bin/docker-compose down || echo 'Compose stop failed'"
Restart=always
RestartSec=10

# Environment variables
Environment=DOCKER_HOST=unix:///var/run/docker.sock

# Resource limits
LimitNOFILE=65536
LimitNPROC=65536

[Install]
WantedBy=multi-user.target
